# salt_root/srv/salt/supervisor/etc/supervisor/conf.d/django_demo.conf
{% set app_port = salt['pillar.get']('config:' + app + ':port', 1234) -%}
{% set app_concurrency = salt['pillar.get']('config:' + app + ':concurrency', 1) -%}
[program:{{ app }}]
command=/opt/{{ app }}/env3/bin/gunicorn -w {{ app_concurrency }} -t 60 --log-level=error --bind 127.0.0.1:{{ app_port }} --max-requests=100 --env DJANGO_SETTINGS_MODULE=root.settings root.wsgi
directory=/opt/{{ app }}/src
user=appuser
group=appuser
autostart=true
autorestart=true
redirect_stderr=True
daemon = False
debug = False
stdout_logfile = /var/log/{{ app }}/gunicorn-out.log
loglevel = "info"


